# packaging/tubee-jobs.service.upstart
description "tubee jobs"
author      "Raffael Sahli <sahli@gyselroth.com>"

start on filesystem or runlevel [2345]
stop on shutdown

script
    exec /usr/bin/tubeecli jobs listen
end script
